# core/crontab
0 3 * * * /usr/bin/docker system prune -f
0 0 * * * /dozer/docker/influxdb/daily_report.sh
30 0 * * * /dozer/docker/paperlessngx/trigger_export.sh
